# sysutils/zmonitor/pkg-plist
bin/zmonitor.py
etc/rc.d/zmonitor
etc/zmonitor.conf.sample
