# sysutils/reflink_dedupe/pkg-plist
bin/reflink_dedupe
etc/rc.d/reflink_dedupe
etc/newsyslog.conf.d/reflink_dedupe.conf
etc/reflink_dedupe.conf.sample
etc/cron.d/reflink_dedupe.cron
